# deploy/perigeos.service
[Unit]
Description=Perigeos - A Kubernetes' Pawn
Documentation=https://github.com/malformed-c/periapsis
After=network-online.target k3s.service
Wants=network-online.target
Requires=dbus.service

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/local/bin/perigeos \
    --perigeosconfig /etc/apsis/perigeos/perigeos.toml \
    --kubeconfig /etc/rancher/k3s/k3s.yaml \
    --base-dir /var/lib/apsis/perigeos \
    --control-socket /run/apsis/perigeos.sock

Restart=on-failure
RestartSec=5s
WatchdogSec=120s

OOMScoreAdjust=-999

# State directory
StateDirectory=apsis
RuntimeDirectory=apsis

# Logging - stdout/stderr go to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=perigeos

# Hardening - perigeos needs root for systemd-nspawn and D-Bus machine1.
# ProtectHome/ProtectSystem that create private mount namespaces must be
# disabled: overlay mounts done here must be visible to transient nspawn
# units that systemd starts in the init mount namespace (PID 1).
ProtectSystem=no
ProtectHome=no
NoNewPrivileges=no
LimitNOFILE=1048576
LimitNPROC=infinity
LimitMEMLOCK=infinity
TasksMax=infinity

# Only kill the main perigeos process - nspawn containers survive both
# stop and restart. HydrateFromRuntime rediscovers running pods on startup.
# For intentional decommission (drain all pods), use: apsis drain
KillMode=process
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
